# example/tollgate.cron
*/10 * * * * root cd /opt/tollgate; ./manage.py refresh_hosts
