# platform/fedora/systemd/tollgate-captivity.service
[Unit]
Description=Tollgate captivity service
After=network.target dbus.service httpd.service

[Service]
PIDFile=/var/run/tollgate-captivity.pid
EnvironmentFile=/etc/sysconfig/tollgate
ExecStart=/usr/sbin/tollgate_captivity -l $TOLLGATEURI

[Install]
WantedBy=multi-user.target
